# agent/poly-shadow.service
[Unit]
Description=Poly-Shadow Polymarket Trading Agent
Documentation=https://github.com/you/poly-shadow-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/poly-shadow-agent/src
EnvironmentFile=/home/ubuntu/poly-shadow-agent/.env
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poly-shadow
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
